# board/mojo2k/raspberryradio/overlay/etc/network/if-up.d/startradio
#!/bin/bash
# Start radio when wifi connection established
/usr/bin/supervisorctl start radio
